# COPS/rootfs/etc/s6-overlay/s6-rc.d/rsync/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

# Start rsync if enabled
if bashio::config.true 'rsync'; then
    bashio::log.green 'Starting rsync daemon...'
    rsync --daemon
fi