# COPS/rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/execlineb -P
# shellcheck shell=bash

# Start PHP-FPM (always required for the web UI)
php-fpm -D -R