# COPS/rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.green "Starting PHP-FPM via s6..."
# -F means "force to stay in foreground"
# -R means "allow running as root" (required for some Alpine setups)
exec php-fpm -F -R