# COPS/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
# shellcheck shell=bash

exec nginx -g "daemon off;" 2>&1