# COPS-S6/rootfs/etc/s6-overlay/s6-rc.d/korrosync/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

# If korrosync is not enabled then exit here
if ! bashio::config.true 'koreader_sync'; then
    exec sleep infinity
fi

# Set Configuration
export KORROSYNC_DB_PATH="/data/korrosync/sync.redb"
export KORROSYNC_SERVER_ADDRESS="127.0.0.1:3000"

bashio::log.info "Starting KORROSYNC..."
exec /usr/bin/korrosync