# COPS/rootfs/etc/s6-overlay/s6-rc.d/cops-init/up
#!/command/with-contenv sh
# shellcheck shell=sh
/etc/s6-overlay/s6-rc.d/cops-init/cops-init.sh